# runtime/claude/telegram/launchd/md.qroky.telegram.listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>md.qroky.telegram.listener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__TG_DIR__/listener.sh</string>
    </array>

    <!-- Polling cadence: one plain-script pass every 30 s. 30 s < the 60 s
         ack criterion (H13) with margin; no ports, no resident process. -->
    <key>StartInterval</key>
    <integer>30</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__TG_DIR__/state/launchd.listener.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__TG_DIR__/state/launchd.listener.stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
